# systemd/borg-backup.timer
[Unit]
Description = Borg backup schedule

[Timer]
OnCalendar = hourly
Persistent = true

[Install]
WantedBy = timers.target
