# systemd/.config/systemd/user/back2spinoza.timer
[Unit]
Description=Execute Borg backup
OnFailure=failure-notify@%n.service

[Timer]
OnBootSec=60
OnUnitActiveSec=3600
Unit=back2spinoza.service

[Install]
WantedBy=multi-user.target
