# systemd/.config/systemd/user/failure-notify@.service
[Unit]
Description=OnFailure notification email for %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/failure-notify.sh
